# deployment/server/crontab
# m h  dom mon dow   command
0 0  * * *  /carlos/deployment/server/scripts/cron_job.sh